# misc/99-power.rules.in
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/bin/su __USER__ -c \"__HOME__/dotfiles/misc/scripts/power-change.sh battery\""
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/bin/su __USER__ -c \"__HOME__/dotfiles/misc/scripts/power-change.sh ac\""
